# crontab
0 18 * * 5 bash /app/run_scrapper.sh > /proc/1/fd/1 2> /proc/1/fd/2
